# templates/singbox-healthcheck.timer
[Unit]
Description=sing-box Health Check Timer
After=sing-box.service

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
