# templates/singbox-healthcheck.timer
[Unit]
Description=Run sing-box health check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
