# templates/singbox-watchdog.service
[Unit]
Description=sing-box Watchdog 2.0 (Go sidecar)
After=sing-box.service
BindsTo=sing-box.service
ConditionPathExists=/usr/local/etc/sing-box/.deployment_config

[Service]
Type=simple
EnvironmentFile=/usr/local/etc/sing-box/.deployment_config
# P1 修复：secret 已从 deployment_config 迁移至 .credentials；
#          通过 LoadCredential 安全注入，Go 代码从 CREDENTIALS_DIRECTORY/dash_secret 读取。
LoadCredential=dash_secret:/usr/local/etc/sing-box/.credentials/dash_secret
# 代理组名和回退节点名与 singbox_build_region_groups.py 保持一致
Environment=PROXY_GROUP_TAG=🚀 节点选择
Environment=WATCHDOG_FALLBACK_NODE=🌏 全局自动
ExecStart=/usr/local/bin/singbox-watchdog
Restart=always
RestartSec=5
KillMode=process
# Jitter handles poll cycles, wait max bounds internally.
WatchdogSec=120

[Install]
WantedBy=multi-user.target
